# crontab
* * * * * /home/pi/send_ibeacon_scan_log 1
